# examples/demo_sharded_snake/deploy/overlays/shard_b/demo-sharded-snake-shard.service
[Unit]
Description=demo_sharded_snake shard_b
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/ssh-system-architect
EnvironmentFile=/srv/ssh-system-architect/examples/demo_sharded_snake/deploy/overlays/shard_b/shard.env
ExecStart=/usr/bin/node /srv/ssh-system-architect/examples/demo_sharded_snake/services/snake_shard/src/main.ts --config /srv/ssh-system-architect/examples/demo_sharded_snake/deploy/overlays/shard_b/shard.config.json
Restart=on-failure
User=snake

[Install]
WantedBy=multi-user.target
